# services/com.nanoclaw.voxtral-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.voxtral-tts</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Update path to your project's .venv/bin/python3 during setup -->
        <string>REPLACE_WITH_VENV_PYTHON</string>
        <string>-m</string>
        <string>mlx_audio.server</string>
        <string>--port</string>
        <string>8771</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/voxtral-tts.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voxtral-tts.err</string>
</dict>
</plist>
